# deployment/scaleway-dev/kubernetes/packer/filesystem/etc/systemd/system/bootstrap.service
[Unit]
Description=Bootstrap
After=vault.service

[Service]
User=root
Type=oneshot
ExecStart=/bootstrap/bootstrap.sh
StandardOutput=journal

[Install]
WantedBy=multi-user.target
